# umh-lite-v2/s6-overlay/services.d/agent/run
#!/bin/sh

# Wait for redpanda to be available
while ! /usr/local/bin/rpk cluster health 2>/dev/null | grep -q "Healthy"; do
  echo "Waiting for Redpanda to be ready..."
  sleep 2
done

exec /usr/local/bin/umh-lite-v2 \
  --config /data/umh-lite-bootstrap.yaml \
  --config-dir /data/configs \
  --log /data/state/agent.log 