# umh-lite-v2/s6-base/s6-rc.d/umh-lite-v2/run
#!/command/execlineb -P
cd /
exec /usr/local/bin/umh-lite-v2 